# contrib/f5gs.service.in
[Unit]
Description=F5 Graceful Scaling

[Service]
ExecStart=@bindir@/f5gs --server
ExecStop=/usr/bin/pkill -9 f5gs
Restart=always

[Install]
WantedBy=multi-user.target
