# contrib/f5gs.socket.in
[Unit]
Description=F5 Graceful Scaling Socket
Conflicts=f5gs.service

[Socket]
Service=f5gs.service
ListenStream=@F5GS_TCP_PORT@
ReusePort=true

[Install]
WantedBy=sockets.target
