# contrib/f5gs.service.in
[Unit]
Description=F5 Graceful Scaling
Documentation=man:f5gs(8)
RequiresMountsFor=@f5gs_rundir@
After=network.target nss-lookup.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=@bindir@/f5gs --server --foreground
Sockets=f5gs.socket
Restart=always
RestartSec=0

[Install]
WantedBy=multi-user.target
