# contrib/f5gs.service.in
[Unit]
Description=F5 Graceful Scaling
After=network.target nss-lookup.target

[Service]
Type=forking
ExecStart=@bindir@/f5gs --server
Restart=always

[Install]
WantedBy=multi-user.target
